# crontab
* * * * * root node /jobs/src/helloworld.js >> /var/log/cron.log 2>&1
